# recipes-graphics/wayland/weston-init/weston.path
[Unit]
Description=Graphic Card Watch
Before=paths.target shutdown.target

[Path]
# DRM backend requires /dev/dri/card0
PathExists=/dev/dri/card0
